# scripts/cicy-code.supervisor.conf
; cicy-code Supervisor Config (Release Mode)
; 使用方式: ln -sf $(pwd)/scripts/cicy-code.supervisor.conf /etc/supervisor/conf.d/cicy-code.conf
;
[program:cicy-code]
command=/usr/local/bin/cicy-code
directory=/home/%(ENV_USER)s
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/cicy-code/error.log
stdout_logfile=/var/log/cicy-code/output.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=HOME="/home/%(ENV_USER)s",TERM="xterm-256color",PATH="/usr/local/bin:/usr/bin:/bin:/home/%(ENV_USER)s/.local/bin:/home/%(ENV_USER)s/.nvm/versions/node/v20.19.0/bin"
