# scripts/cicy-code.supervisor.conf
; cicy-code Supervisor Config (Release Mode)
; 使用方式: ln -sf $(pwd)/scripts/cicy-code.supervisor.conf /etc/supervisor/conf.d/cicy-code.conf
;
[program:cicy-api]
command=/bin/bash %(ENV_HOME)s/projects/cicy-code/dev.sh
directory=%(ENV_HOME)s/projects/cicy-code
user=w3c_offical
autostart=true
autorestart=true
stdout_logfile=/var/log/cicy-api.log
stderr_logfile=/var/log/cicy-api.log