# scripts/cicy-code-dev.supervisor.conf
; cicy-code Supervisor Config (Development Mode)
; 使用方式: ln -sf $(pwd)/scripts/cicy-code-dev.supervisor.conf /etc/supervisor/conf.d/cicy-code.conf
;
; 注意: directory 必须指向 api/ 目录，这样 --dev 模式才能从 mgr/resources/ 读取 inject HTML
;
[program:cicy-code]
command=%(ENV_HOME)s/projects/cicy-code/api/mgr --dev
directory=%(ENV_HOME)s/projects/cicy-code/api
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/cicy-code/error.log
stdout_logfile=/var/log/cicy-code/output.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=HOME="/home/%(ENV_USER)s",TERM="xterm-256color",PATH="/usr/local/bin:/usr/bin:/bin:/home/%(ENV_USER)s/.local/bin:/home/%(ENV_USER)s/.nvm/versions/node/v20.19.0/bin"
