# scripts/code-server.supervisor.conf
; code-server Supervisor Config (独立部署时使用)
; 注意: cicy-code 默认自动管理 code-server，此配置仅用于需要单独管理 code-server 的场景
; 使用方式: ln -sf $(pwd)/scripts/code-server.supervisor.conf /etc/supervisor/conf.d/code-server.conf
;
[program:code-server]
command=code-server --bind-addr 127.0.0.1:18080 --auth none /home/%(ENV_USER)s
directory=/home/%(ENV_USER)s
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/cicy-code/code-server-error.log
stdout_logfile=/var/log/cicy-code/code-server-output.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
environment=HOME="/home/%(ENV_USER)s",TERM="xterm-256color"
