# docker/cloudrun/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:mysql]
command=/usr/sbin/mysqld --user=mysql --datadir=/var/lib/mysql
autostart=true
autorestart=true

[program:redis]
command=redis-server --appendonly no --save ""
autostart=true
autorestart=true

[program:code-server]
command=code-server --bind-addr 0.0.0.0:8002 --auth none /home/w3c_offical
user=w3c_offical
environment=HOME="/home/w3c_offical",PORT="8002"
autostart=true
autorestart=true
stdout_logfile=/var/log/code-server.log
stderr_logfile=/var/log/code-server.log

[program:cicy-api]
command=/bin/bash /home/w3c_offical/projects/cicy-code/api/run.sh
directory=/home/w3c_offical/projects/cicy-code/api
user=w3c_offical
environment=HOME="/home/w3c_offical",MYSQL_DSN="root:@tcp(127.0.0.1:3306)/cicy_code?parseTime=true",REDIS_ADDR="127.0.0.1:6379"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
